# kasa-bridge.service

# /etc/systemd/system/kasa-bridge.service
[Unit]
Description=Kasa MQTT Bridge Service
After=network-online.target
Wants=network-online.target
StartLimitIntervalSec=60
StartLimitBurst=10

[Service]
Type=notify
NotifyAccess=main
User=<user>
WorkingDirectory=/home/<user>/kasa_mqtt_simple_bridge
Environment=PYTHONUNBUFFERED=1
ExecStart=/home/<user>/kasa_mqtt_simple_bridge/venv/bin/python -u /home/<user>/kasa_mqtt_simple_bridge/kasa_>

# Policy di riavvio
Restart=on-failure
RestartSec=5s

# Watchdog: il codice invia WATCHDOG=1 ogni ~15s
WatchdogSec=120s

# (opzionale) limita l’attesa iniziale
TimeoutStartSec=15s

StandardOutput=journal
StandardError=journal

# (opzionali)
# TimeoutStopSec=15s
# KillMode=process

[Install]
WantedBy=multi-user.target
